# container-files/etc/supervisor.d/jira.conf
[program:jira]
command = /opt/atlassian/jira/bin/start-jira.sh -fg
autorestart = true
